# 手把手带你自制Linux系统之二 简易Linux的制作/inittab


id:3:initdefault:
si::sysinit:/etc/rc.d/rc.sysinit
