# 手把手带你自制Linux系统之六 编译内核及busybox完成系统定制/inittab
::sysinit:/etc/rc.d/rc.sysinit
::respawn:/sbin/getty 9600 tty1
::respawn:/sbin/getty 9600 tty2
::shutdown:/bin/umount -a -r
::ctrlaltdel:/sbin/reboot
